# TestSuite.service
[Unit]
Description=TestSuite Service
After=network.target

[Service]
WorkingDirectory=/usr/local/bin/TestSuite
ExecStart=/usr/local/bin/TestSuite/TestSuite
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
